# scripts/crontab
0 0 1 * * python /home/getWinner.py >> /logfile